# odoPDS/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod
priority=0
startsecs=5

[program:openpds]
command=/bin/bash -c "/start.sh"
priority=10
startsecs=10


